# .config/systemd/user/battery-logger.service
[Unit]
Description=Log battery discharge rate to history.json
After=basic.target

[Service]
Type=oneshot
ExecStart=%h/.config/quickshell/battery_popup/log_battery
